# deployment/systemd/ohana-agent.service
[Unit]
Description=Ohana-Agent infrastructure observation service
Documentation=https://github.com/cedric-HAOS/Ohana-Agent
Wants=network-online.target
After=network-online.target

[Service]
Type=simple

User=ohana-agent
Group=ohana-agent

ExecStart=/opt/ohana-agent/venv/bin/ohana-agent \
    --config /etc/ohana-agent/shikamaru.yaml \
    --infrastructure /etc/ohana-agent/infrastructure.yaml \
    --dns-config /etc/ohana-agent/plugins/dns.yaml \
    --log-level INFO

Restart=on-failure
RestartSec=5s

TimeoutStopSec=30s
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ohana-agent

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target